# scripts/kalshi-mcp-tunnel.service
# =============================================================================
# kalshi-mcp-tunnel.service — systemd unit template
#
# THIS IS A TEMPLATE — DO NOT COPY DIRECTLY. Use scripts/install_tunnel_service.sh
# to install this service, which will stamp the correct paths for your machine.
#
# The install script generates /etc/systemd/system/kalshi-mcp-tunnel.service
# with the correct User= and WorkingDirectory= for this machine.
# =============================================================================

[Unit]
Description=Kalshi MCP Server + GCE Reverse SSH Tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Load deployment config from .env
EnvironmentFile=-__PROJECT_DIR__/.env

# GCP Tunnel Configuration (from .env or Environment= in override)
Environment=GCP_PROJECT=${GCP_PROJECT}
Environment=VM_INSTANCE=${VM_INSTANCE}
Environment=ZONE=${ZONE:-us-central1-a}
Environment=MCP_PORT=${MCP_PORT:-8765}

# Machine-specific paths — stamped by install_tunnel_service.sh
User=__USER__
WorkingDirectory=__PROJECT_DIR__

ExecStart=/bin/bash __PROJECT_DIR__/scripts/tunnel_mcp.sh \
    ${GCP_PROJECT} \
    ${VM_INSTANCE} \
    ${ZONE}

Restart=on-failure
RestartSec=20s

TimeoutStartSec=120s
TimeoutStopSec=90s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=kalshi-mcp-tunnel

[Install]
WantedBy=multi-user.target